# postgresql/templates/default/pgpool.upstart.erb
description "pgpool daemon"

start on runlevel [2]
stop on runlevel [016]

console owner

exec /usr/bin/pgpool -n -d -f /etc/pgpool.conf > /var/log/postgresql/pgpool.log 2>&1

post-stop exec /usr/bin/pgpool -m immediate stop

respawn
